# scripts/systemd/onex-tick-merge-sweep.service
[Unit]
Description=ONEX Merge-Sweep Tick — arm auto-merge, update-branch, rerun stale gates
Documentation=https://github.com/OmniNode-ai/omnibase_infra
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.omnibase/.env
ExecStart=/bin/bash %h/Code/omni_home/omnibase_infra/scripts/ticks/merge-sweep-tick.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=onex-tick-merge-sweep

[Install]
WantedBy=default.target
