# scripts/deploy-agent/deploy/deploy-agent.service
# SPDX-FileCopyrightText: 2025 OmniNode.ai Inc.
# SPDX-License-Identifier: MIT
#
# User-scope systemd unit for the ONEX deploy agent.
# Install: cp deploy-agent.service ~/.config/systemd/user/
#          systemctl --user daemon-reload
#          systemctl --user enable --now deploy-agent.service
#
# WatchdogSec: agent must call sd_notify("WATCHDOG=1") at least once per 60s.
# Type=notify: systemd waits for READY=1 before considering startup complete.
# Restart=on-failure: restart on non-zero exit; not on clean SIGTERM.

[Unit]
Description=ONEX Deploy Agent — Kafka-driven rebuild orchestrator
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=notify
WorkingDirectory=/data/omninode/deploy-agent
EnvironmentFile=/home/jonah/.omnibase/.env
Environment=DEPLOY_AGENT_STATE_DIR=/data/omninode/deploy-agent/state/jobs
Environment=DEPLOY_AGENT_PORT=8099
ExecStartPre=/bin/sh -c 'fuser -k 8099/tcp || true'
ExecStart=/data/omninode/deploy-agent/venv/bin/python -m deploy_agent
Restart=on-failure
RestartSec=5
WatchdogSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deploy-agent

[Install]
WantedBy=default.target
