# scripts/deploy-agent/deploy/deploy-agent.service
# SPDX-FileCopyrightText: 2025 OmniNode.ai Inc.
# SPDX-License-Identifier: MIT
#
# User-scope systemd unit for the ONEX deploy agent.
# Install: cp deploy-agent.service ~/.config/systemd/user/
#          systemctl --user daemon-reload
#          systemctl --user enable --now deploy-agent.service
#
# No WatchdogSec (OMN-13760): this agent's core job is orchestrating docker
#   rebuilds that run as synchronous, minutes-long subprocess.run() calls inside
#   the asyncio event loop (see deploy_agent.executor). A liveness watchdog
#   cannot distinguish a healthy in-progress rebuild from a genuine hang — any
#   ping cadence loose enough to let a legitimate multi-minute rebuild survive is
#   also loose enough to mask a real hang, so the watchdog adds risk without
#   detection value. Restart=on-failure already recovers genuine crashes
#   (non-zero exit). A previous WatchdogSec=30 SIGABRTed the process every 30s
#   because the ping depended on the optional systemd-python package, which was
#   absent from the deployed venv, so sd_notify("WATCHDOG=1") silently no-op'd.
# Type=simple: startup is considered complete immediately after ExecStart forks.
#   We do NOT use Type=notify because systemd-python may not be installed in the
#   venv, which would cause a Type=notify deadlock (READY=1 never sent → killed
#   after TimeoutStartSec).
# Restart=on-failure: restart on non-zero exit; not on clean SIGTERM.

[Unit]
Description=ONEX Deploy Agent — Kafka-driven rebuild orchestrator
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/data/omninode/omnibase_infra/scripts/deploy-agent
EnvironmentFile=/home/jonah/.omnibase/.env
Environment=DEPLOY_AGENT_DIR=/data/omninode/omnibase_infra/scripts/deploy-agent
Environment=DEPLOY_AGENT_STATE_DIR=/data/omninode/deploy-agent/state/jobs
Environment=DEPLOY_AGENT_PORT=8099
ExecStartPre=/bin/sh -c 'fuser -k 8099/tcp || true'
ExecStart=/data/omninode/omnibase_infra/scripts/deploy-agent/.venv/bin/python -m deploy_agent
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deploy-agent

[Install]
WantedBy=default.target
