# deploy/systemd/onex-log-monitor.service
[Unit]
Description=OmniNode Log Monitor — Docker containers + Slack alerts
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/home/jonah/.omnibase/.env
Environment="MONITOR_PROJECTS=omnibase-infra,omnibase-infra-runtime"
Environment="MONITOR_FILE_LOGS=/home/jonah/.onex_state/logs/env-sync.log:/home/jonah/.onex_state/logs/hooks.log:/home/jonah/.onex_state/logs/pipeline-trace.log"
Environment="MONITOR_JOURNALS=deploy-agent.service"
ExecStart=/usr/bin/python3 /home/jonah/.omnibase/infra/deployed/0.34.1/scripts/monitor_logs.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
