# scripts/deploy-agent/deploy/deploy-agent.service.d/override.conf
# SPDX-FileCopyrightText: 2025 OmniNode.ai Inc.
# SPDX-License-Identifier: MIT
#
# systemd drop-in override for deploy-agent.service on the .201 host
# (omninode-pc). Install: copy to
#   ~/.config/systemd/user/deploy-agent.service.d/override.conf
# then `systemctl --user daemon-reload && systemctl --user restart deploy-agent`.
#
# Why this file is tracked (OMN-13760): the live .201 override was previously an
# untracked, hand-edited drop-in — a durable-evidence gap. This captures the
# deployed state so the running unit matches a source-of-truth file.
#
# WatchdogSec=0 explicitly disables the systemd liveness watchdog even if an
# older base unit that still carries WatchdogSec=30 is installed. See
# deploy/deploy-agent.service for the full rationale (minutes-long synchronous
# rebuilds block the event loop; a liveness watchdog cannot distinguish a healthy
# in-progress rebuild from a hang, and Restart=on-failure recovers real crashes).
#
# NOTE (residual drift, tracked under OMN-13810): the ExecStart below runs code
# from the omni_home canonical clone (via PYTHONPATH) using the legacy standalone
# venv at /data/omninode/deploy-agent/venv. Reconciling the three divergent code
# copies + venv onto a single canonical path is OMN-13810, out of scope here.

[Service]
WatchdogSec=0
WorkingDirectory=/data/omninode/omni_home/omnibase_infra/scripts/deploy-agent
Environment=PATH=/home/jonah/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/data/omninode/omni_home/omnibase_infra/scripts/deploy-agent
Environment=KAFKA_ENVIRONMENT=stability-test
ExecStart=
ExecStart=/usr/bin/env PATH=/home/jonah/.local/bin:/usr/local/bin:/usr/bin:/bin PYTHONPATH=/data/omninode/omni_home/omnibase_infra/scripts/deploy-agent KAFKA_BOOTSTRAP_SERVERS=127.0.0.1:39092 KAFKA_ENVIRONMENT=stability-test /data/omninode/deploy-agent/venv/bin/python -m deploy_agent
